# apps/pro/ui/config/build/common/templates/linux/units/metasploit-postgresql-env.service.erb
[Unit]
Description=Metasploit Pro PostgreSQL environment generator
After=network.target

[Service]
Type=oneshot
RemainAfterExit=true
User=root
Group=root
ExecStart=@@BITROCK_INSTALLDIR@@/postgresql/scripts/gen_env.sh @@BITROCK_INSTALLDIR@@/postgresql/scripts/.env.txt
