# apps/pro/ui/config/build/common/templates/linux/units/metasploit.target.erb
[Unit]
Description=Metasploit Pro services
After=network.target

Wants=metasploit.service metasploit-prosvc.service metasploit-ui.service metasploit-worker.service
PartOf=metasploit.service

[Install]
WantedBy=multi-user.target
